# revsw-proxy-config/configs/crontab/wallarm
#
# cron job for Wallarm WAF
#

PATH="/usr/share/wallarm-common:/bin:/sbin:/usr/bin:/usr/sbin"

*/15 *  * * *   root      timeout -k 15s 14m /usr/share/wallarm-common/syncnode -r 600

